# begemot.conf
[program: begemot]
directory=/home/begemot
command=python app.py

autostart=true
autorestart=true
stderr_logfile=/home/begemot/logfile_err.log
stdout_logfile=/home/begemot/logfile.log
