# point.conf
[program: point]
directory=/home/point
command=python app.py

autostart=true
autorestart=true
stderr_logfile=/home/point/logfile_err.log
stdout_logfile=/home/point/logfile.log
