# user_bot.conf
[program: user]
directory=/home/user_vendista
command=python app.py

autostart=true
autorestart=true
stderr_logfile=/home/user_vendista/logfile_err.log
stdout_logfile=/home/user_vendista/logfile.log
